# .config/systemd/user/battery-check.service
[Unit]
Description=Check battery status

[Service]
Type=oneshot
ExecStart=%h/.local/lib/timers/battery-warn
